# archiso_profile/add/airootfs/etc/skel/.config/autostart/os-installer.desktop
[Desktop Entry]
Version=1.0
Type=Application
Name=Startup
Comment=Start the OS installer
TryExec=os-installer
Exec=os-installer
